# files/etc/uci-defaults/99-printer
#!/bin/sh

/etc/init.d/p910nd enable
/etc/init.d/hp1020 enable

exit 0
